# launchd/agents/com.joshuashew.xlsx-clip-watcher.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joshuashew.xlsx-clip-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>DOTFILES_DIR/launchd/scripts/xlsx-clip-watcher.sh</string>
    </array>
    <!-- WatchPaths: launchd watches ~/Downloads natively; no fswatch needed -->
    <key>WatchPaths</key>
    <array>
        <string>HOME_DIR/Downloads</string>
    </array>
    <!-- Run once at load to catch any file that landed before the agent started -->
    <key>RunAtLoad</key>
    <true/>
    <!-- Throttle: don't re-trigger more than once per 5 seconds -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <!-- Log to state dir (created by install script before plist is loaded) -->
    <key>StandardOutPath</key>
    <string>HOME_DIR/.local/state/xlsx-clip-watcher/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_DIR/.local/state/xlsx-clip-watcher/launchd.log</string>
</dict>
</plist>
