# launchd/agents/com.joshuashew.xlsx-clip-watcher.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joshuashew.xlsx-clip-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>DOTFILES_DIR/launchd/scripts/xlsx-clip-watcher.sh</string>
    </array>
    <!-- Poll every 5 s. WatchPaths was dropped: macOS 12+ does not fire
         WatchPaths reliably for browser downloads (Chrome/Safari write via
         atomic temp-file rename; launchd fires on the crdownload/tmp appear,
         not on the final xlsx rename, so every scan ran before the file
         existed and logged "done (0 processed)"). StartInterval is slower
         but 100% reliable. RunAtLoad catches files already present at load. -->
    <key>StartInterval</key>
    <integer>5</integer>
    <key>RunAtLoad</key>
    <true/>
    <!-- Log to state dir (created by install script before plist is loaded) -->
    <key>StandardOutPath</key>
    <string>HOME_DIR/.local/state/xlsx-clip-watcher/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>HOME_DIR/.local/state/xlsx-clip-watcher/launchd.log</string>
</dict>
</plist>
